# repocket/rootfs/etc/services.d/repocket/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Run the Repocket service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================
export RP_EMAIL=$(bashio::config 'email')
export RP_API_KEY=$(bashio::config 'key')

bashio::log.info "Starting Repocket Service..."

## Call the app
node /app/dist/index.js
