# bitping/rootfs/etc/services.d/bitping/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Run the Bitping service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================
declare -a parameters

bashio::log.info "Starting Bitping Service..."

parameters+=(--server)

bashio::log.debug "Arch: $(uname -m)"

if [[ $(uname -m) != "armv7l" ]]; then
  /opt/bitping/release/bitping-node-amd64-linux "${parameters[@]}"
else
  /opt/bitping/release/bitping-node-armv7-linux "${parameters[@]}"
fi
