# bitping/rootfs/etc/services.d/bitping/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Run the Bitping service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Starting Bitping Service..."

if [[ "$(bashio::config 'email')" == "" || "$(bashio::config 'password')" == "" ]]; then
  bashio::exit.nok "Bitping email and/or password are empty."
fi

export BITPING_EMAIL="$(bashio::config 'email')"
export BITPING_PASSWORD="$(bashio::config 'password')"

/opt/bitping/bitpingd